# config/common/files/etc/uci-defaults/9999-dezeroku-startup
chmod +x /etc/init.d/dezeroku_startup
/etc/init.d/dezeroku_startup enable
