# config/common/files/etc/uci-defaults/9999-openwrt_build_wrapper_startup
chmod +x /etc/init.d/openwrt_build_wrapper_startup
/etc/init.d/openwrt_build_wrapper_startup enable
